# contrib/init/bitcoin5000d.conf
description "Bitcoin5000 Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BITCOIN5000D_BIN="/usr/bin/bitcoin5000d"
env BITCOIN5000D_USER="bitcoin5000"
env BITCOIN5000D_GROUP="bitcoin5000"
env BITCOIN5000D_PIDDIR="/var/run/bitcoin5000d"
# upstart can't handle variables constructed with other variables
env BITCOIN5000D_PIDFILE="/var/run/bitcoin5000d/bitcoin5000d.pid"
env BITCOIN5000D_CONFIGFILE="/etc/bitcoin5000/bitcoin5000.conf"
env BITCOIN5000D_DATADIR="/var/lib/bitcoin5000d"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # bitcoin5000d will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BITCOIN5000D_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitcoin5000d."
        echo "The setting must appear in $BITCOIN5000D_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin5000 Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BITCOIN5000D_PIDDIR"
    chmod 0755 "$BITCOIN5000D_PIDDIR"
    chown $BITCOIN5000D_USER:$BITCOIN5000D_GROUP "$BITCOIN5000D_PIDDIR"
    chown $BITCOIN5000D_USER:$BITCOIN5000D_GROUP "$BITCOIN5000D_CONFIGFILE"
    chmod 0660 "$BITCOIN5000D_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BITCOIN5000D_PIDFILE" \
    --chuid $BITCOIN5000D_USER:$BITCOIN5000D_GROUP \
    --exec "$BITCOIN5000D_BIN" \
    -- \
    -pid="$BITCOIN5000D_PIDFILE" \
    -conf="$BITCOIN5000D_CONFIGFILE" \
    -datadir="$BITCOIN5000D_DATADIR" \
    -disablewallet \
    -daemon

